# templates/default/redis-slave.erb
#!/bin/sh
#
# Simple Redis init.d script conceived to work on Linux systems
# as it does use of the /proc filesystem.
#
# description: Redis is an in memory key-value store database
#
### BEGIN INIT INFO
# Provides: redis-slave
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Required-Start: 
# Required-Stop: 
# Description: redis-slave init script
### END INIT INFO

REDISNAME=slave
REDISPORT=6380
EXEC="cd <%= node.override['redisio']['install_dir'] %>;<%= node.override['redisio']['install_dir'] %>/src/redis-server <%= node.override['redisio']['install_dir'] %>/redis-slave.conf --port <%= @slaveport %> --cluster-enabled yes --cluster-config-file <%= node.override['redisio']['install_dir'] %>/node-<%= @slaveport %>.conf --cluster-node-timeout 10000 --appendonly yes --appendfilename appendonly-<%= @slaveport %>.aof --dbfilename dump-<%= @slaveport %>.rdb --logfile <%= node.override['redisio']['install_dir'] %>/<%= @slaveport %>.log --daemonize yes --protected-mode no"
CLIEXEC=/redis-setup/redis-3.2.5/src/redis-cli


PIDFILE=/var/run/redis-slave.pid

ulimit -n 10032

case "$1" in
    status)
        if [ -f $PIDFILE ]
        then
          echo "redis-$REDISNAME $PIDFILE exists, pid is $(cat $PIDFILE), should be running"
          ps -p $(cat $PIDFILE) >/dev/null 2>&1
          exit $?
        else
          echo "redis-$REDISNAME $PIDFILE doesn't exist"
          exit 3
        fi
        ;;
    start)
        if [ -f $PIDFILE ]
        then
                echo "$PIDFILE exists, process is already running or crashed"
        else
                echo "Starting Redis server..."
                eval $EXEC
        fi
        ;;
    stop)
        if [ ! -f $PIDFILE ]
        then
                echo "$PIDFILE does not exist, process is not running"
        else
                PID=$(cat $PIDFILE)
                echo "Stopping ..."

                
                $CLIEXEC  -p 6380 shutdown

                while [ -x /proc/${PID} ]
                do
                    echo "Waiting for Redis to shutdown ..."
                    sleep 1
                done
                echo "Redis stopped"
        fi
        ;;
    *)
        echo "Please use start or stop as first argument"
        ;;
esac
